# scripts/midimind.service
[Unit]
Description=MidiMind - Professional MIDI Orchestration System
After=network.target sound.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/midimind
ExecStart=/opt/midimind/bin/midimind --config /etc/midimind/config.json
Restart=always
RestartSec=3

# Permissions
SupplementaryGroups=audio dialout plugdev

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=midimind

[Install]
WantedBy=multi-user.target